# supervisord.conf
[supervisord]
nodaemon=true

[program:react]
command=npm start

[program:flask]
command=python3 flask/app.py

