# supervisord.conf
[supervisord]
nodaemon=true

[program:react]
command=yarn start

[program:flask]
command=yarn start-flask

