# supervisord.conf
[supervisord]
nodaemon=true

[program:react]
command=npm start

[program:flask]
command=npm run start-api
