# supervisord.conf
[supervisord]
nodaemon=true

[program:react]
command=yarn start

[program:flask]
command=python3 flask/app.py

